# compose-update@.service
[Unit]
Description=docker-compose update service "%f"

[Service]
Type=oneshot
ExecStart=/usr/bin/env systemctl reload compose@%i.service
