# compose-update@.timer
[Unit]
Description=docker-compose update timer "%f"
PartOf=compose@.service

[Timer]
OnCalendar=*-*-* 00:02:00
RandomizedDelaySec=1h
