# docker-image-prune.service
[Service]
Type=oneshot
ExecStart=/usr/bin/env docker image prune -af
